# service/system/etc/cron.daily/update_check
#!/bin/sh
sudo apt-get update && apt-get -s -o Debug::NoLocking=true upgrade | grep ^Inst |\
 cut -f '2-4' -d ' ' | tr -d '[' | tr -d ']' | tr -d '(' >/usr/lib/piwebagent2/update_check 
chown piwebagent2 /usr/lib/piwebagent2/update_check